# demo-routing/load-balancer/services/consul_haproxy.ini
[program:consul_haproxy]
command=/usr/local/bin/consul-template -consul consul:8500 -template "/etc/haproxy/haproxy.cfg.ctmpl:/etc/haproxy/haproxy.cfg:supervisorctl reload haproxy"
stdout_logfile=syslog
stderr_logfile=syslog
