# demo-common/base/services/consul.ini
[program:consul]
command=/usr/local/bin/consul agent -data-dir /data/consul -join=consul -config-dir /etc/consul.d 2>&1
