# demo-common/base/services/dnsmasq.ini
[program:dnsmasq]
command=/usr/local/sbin/dnsmasq_daemon.sh 2>&1
