# demo-common/base/services/syslog-ng.ini
[program:syslog-ng]
command=/usr/sbin/syslog-ng -f /etc/syslog-ng/syslog-ng.conf --foreground
stdout_logfile=syslog
stderr_logfile=syslog
