# demo-routing/proxy/services/consul_nginx.ini
[program:consul_nginx]
command=/usr/local/bin/consul-template -consul consul:8500 -template "/etc/nginx/sites-enabled/proxy.conf.ctmpl:/etc/nginx/sites-enabled/proxy.conf:supervisorctl reload nginx"
stdout_logfile=syslog
stderr_logfile=syslog
